# deploy/cron/ts6-refresh.cron
0 11 * * * root cd /opt/ts6-manager && /bin/bash scripts/refresh-ts-server.sh >> /var/log/ts6-refresh.log 2>&1
